# callattendant/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: SQLite Web
# Runs the NGINX daemon
# ==============================================================================

# Wait for SQLLite Web to start
bashio::net.wait_for 5000

bashio::log.info "Starting NGINX..."

chown -R nginx:nginx /run/nginx /var/lib/nginx

# Start NGINX
exec nginx -c /etc/nginx/nginx.conf -g "daemon off;"
